# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=1

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autorestart=true
priority=2
